# CTIN/Spacewalk/Spacewalk/supervisor.conf
[supervisord]
nodaemon=true

[program:tomcat]
command=/usr/libexec/tomcat/server start

[program:httpd]
command=/usr/sbin/httpd -D FOREGROUND

[program:cobblerd]
command=/usr/bin/cobblerd -F

